# ml_trade_btc_1m.service
[Unit]
Description=ML Trading Bot — BTC/USDT 1m
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=luiz
WorkingDirectory=/home/luiz/ml_trade
ExecStart=/usr/bin/python3 /home/luiz/ml_trade/main.py --loop --config config_btc_1m.yaml
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/luiz/ml_trade/logs/bot_btc_1m.log
StandardError=append:/home/luiz/ml_trade/logs/bot_btc_1m.log

StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
